# crawling/com.sangsangfinder.crawler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sangsangfinder.crawler</string>

    <!--
      ProgramArguments: 첫 번째 항목이 실행 파일, 나머지가 인수
      .venv311/bin/python 이 없으면 `which python` 결과로 교체
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mac/Projects/sangsangfinder/.venv311/bin/python</string>
        <string>/Users/mac/Projects/sangsangfinder/crawling/auto_crawler.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mac/Projects/sangsangfinder</string>

    <!-- 3600초 = 1시간 간격 실행 -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <!-- 로드 즉시 1회 실행 -->
    <key>RunAtLoad</key>
    <true/>

    <!-- stdout/stderr 를 단일 파일로 합산 -->
    <key>StandardOutPath</key>
    <string>/Users/mac/Projects/sangsangfinder/crawling/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mac/Projects/sangsangfinder/crawling/launchd.log</string>

    <!-- 비정상 종료 시 자동 재시작 하지 않음 (cron처럼 스케줄 기반으로만 실행) -->
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
